# idl/gdas/crontab/update_field_data.sh
#!/bin/bash

/home/gsbowers/github/hawc/idl/lightning_mon/get_field_data.py
